# systemd/api-maker-scheduler.service
[Unit]
Description=API Maker Scheduler
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/api_maker
Environment=HOME=/root
Environment=PATH=/root/api_maker/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=DISPLAY=:1
ExecStart=/bin/bash -o pipefail -c 'mkdir -p /root/api_maker/logs && exec /usr/bin/python3 /root/api_maker/scripts/scheduler.py 2>&1 | /usr/bin/tee -a /root/api_maker/logs/keyhive-scanner.log'
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
