# systemd/keyhive-web.service
[Unit]
# Serves the Web UI that reports scanner/proxy state and exposes the local
# control endpoints behind FastAPI.
Description=KeyHive Web UI
After=network-online.target
Wants=network-online.target

[Service]
# The UI binds publicly by default for operator access; firewall it or change
# the host before treating this as anything beyond a private control surface.
Type=simple
User=root
WorkingDirectory=@PROJECT_DIR@
Environment=HOME=/root
Environment=PATH=@VENV_BIN@:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=KEYHIVE_WEB_HOST=0.0.0.0
Environment=KEYHIVE_WEB_PORT=8080
ExecStart=@PYTHON_BIN@ -m uvicorn web_ui.back_end.app:app --host 0.0.0.0 --port 8080
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
