# systemd/keyhive-web.service
[Unit]
# Serves the Web UI that reports scanner/proxy state and exposes the local
# control endpoints behind FastAPI.
Description=KeyHive Web UI
After=network-online.target
Wants=network-online.target

[Service]
# The UI stays on localhost by default; the host/port values are still exposed
# so the installer and dashboard can report them consistently.
Type=simple
User=root
WorkingDirectory=/root/api_maker
Environment=HOME=/root
Environment=PATH=/root/api_maker/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=KEYHIVE_WEB_HOST=0.0.0.0
Environment=KEYHIVE_WEB_PORT=8080
ExecStart=/root/api_maker/.venv/bin/python -m uvicorn web_ui.back_end.app:app --host 0.0.0.0 --port 8080
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
